# com.jp.weekly-preview-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jp.weekly-preview-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jp/Desktop/projects/weekly-preview-assistant/venv/bin/python</string>
        <string>/Users/jp/Desktop/projects/weekly-preview-assistant/main.py</string>
        <string>--check</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jp/Desktop/projects/weekly-preview-assistant</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jp/Desktop/projects/weekly-preview-assistant/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jp/Desktop/projects/weekly-preview-assistant/logs/launchd-stderr.log</string>
</dict>
</plist>
